# crontabs
0 12 * * * screen -S tf2 -X stuff 'sm_say Changing Map to Rainytower in 1 Hour!^M' 30 12 * * * screen -S tf2 -X stuff 'sm_say Changing Map to Rainytower in 30 Minute>50 12 * * * screen -S tf2 -X stuff 'sm_say Changing Map to Rainytower in 10 Minute>59 12 * * * screen -S tf2 -X stuff 'sm_say Changing Map to Rainytower in 1 Minute!>0 13 * * * screen -S tf2 -X stuff 'changelevel plr_hightower_rainy_bum^M'

0 0 * * * screen -S tf2 -X stuff 'sm_say Changing Map to Rainytower in 1 Hour!^M'
30 0 * * * screen -S tf2 -X stuff 'sm_say Changing Map to Rainytower in 30 Minutes>50 0 * * * screen -S tf2 -X stuff 'sm_say Changing Map to Rainytower in 10 Minutes>59 0 * * * screen -S tf2 -X stuff 'sm_say Changing Map to Rainytower in 1 Minute!^>0 1 * * * screen -S tf2 -X stuff 'changelevel plr_hightower_rainy_bum^M'

@reboot cd /hlserver && screen -dmS tf2 && screen -S tf2 -X stuff './tf.sh^M'
